# zenpool-hub.service
[Unit]
Description=ZenPool Hub — distributed key proxy for OpenCode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/zenpool/zenpool.py hub
WorkingDirectory=/opt/zenpool
Restart=always
RestartSec=5
User=root
Environment=ZENPOOL_PORT=5051
Environment=ZENPOOL_DATA=/opt/zenpool/zenpool-data.json
Environment=ZENPOOL_SELF=/opt/zenpool/zenpool.py

[Install]
WantedBy=multi-user.target
